# homework_2/playbook/templates/vector.service.j2
[Unit]
Description=Vector Log Aggregator
After=network-online.target

[Service]
ExecStart=/usr/local/bin/vector --config /etc/vector.yml
Restart=always

[Install]
WantedBy=multi-user.target